# ansible-roles/ubuntu/files/etc/profile.d/bash-shared-history.sh
#!\bin\sh

HISTSIZE=5000
HISTFILESIZE=5000
shopt -s histappend
