# ansible-roles/bridged-networking/files/etc/network/if-up.d/make-bridged-default-route
#!/bin/bash

GATEWAY_IP=$(cat /var/lib/dhcp/*eth1.leases | grep -oP "routers \K[^;]+" | head -1)

if [[ !  -z  $GATEWAY_IP  ]]; then
    ip route change to default dev eth1 via $GATEWAY_IP
fi
